# deploy/backup/runspace-db-backup.timer
[Unit]
Description=Daily RunSpace database backup

[Timer]
OnCalendar=*-*-* 03:17:00 UTC
RandomizedDelaySec=10m
Persistent=true
Unit=runspace-db-backup.service

[Install]
WantedBy=timers.target
